# mariadb/os/etc/profile.d/mariadb.sh
export PATH=/usr/local/mariadb/bin:$PATH
