# iptable/iptable-init.service
[Unit]
After=network.target
Description=iptable-init

[Service]
Type=oneshot
ExecStart=DIR/init.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target


